# lib/cookbooks/toque/templates/default/resque-X.conf.erb
start on starting <%= @master_config %>
stop on stopping <%= @master_config %>

respawn

description "Resque Worker"

# run the thing:
exec su - <%= @deploy_user %> -c 'cd "<%= @current_path %>"; export RAILS_ENV=<%= @rails_env %>; export QUEUE="<%= @queues %>"; bundle exec rake resque:work >> "<%= @shared_path %>/log/resque-<%= @worker_number %>.log" 2>&1'
